# service/privoxy/run
#!/bin/sh
chpst -u privoxy privoxy --no-daemon
